# mox/files/rc.local
iw dev mlan0 interface add uap0 type __ap

